# supervisor-app.conf
[supervisord]

[program:app-uwsgi]
command = /usr/local/bin/uwsgi  --socket app.sock --wsgi-file /rapidpro/temba/wsgi.py --chmod-socket=666

[program:nginx-app]
command = /usr/sbin/nginx
